# log_sender.service
[Unit]
Description=Security Service
After=network.target

[Service]
ExecStart=/path/to/install_security.sh
WorkingDirectory=/opt/security
Restart=always
RestartSec=5
User=your_username
Group=your_group

[Install]
WantedBy=multi-user.target
